# radxa/radxa-webserver.service
[Unit]
Description=Radxa Web Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/radxa/www/server.py
WorkingDirectory=/home/radxa/www
Restart=always
RestartSec=5
Environment=PORT=4000

[Install]
WantedBy=default.target
